# linux-scheduler/daily-checkin.service
[Unit]
Description=Daily Check-in Service - Wake Windows PC and Execute Check-in Tasks
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/satellite-y/daily_tasks.sh
User=dietpi
Group=dietpi

# 环境变量配置（可选）
# Environment="COMET_API_KEY=your-api-key-here"

# 日志
StandardOutput=journal
StandardError=journal

# 超时设置（5分钟应该足够）
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
